# trunk/wikipathways/wpi/crontab
00 00 * * * /var/www/wikipathways/wpi/sitemap/runSiteMap.sh
00 00 * * * rm /var/www/wikipathways/wpi/tmp/*
